# files/etc/uci-defaults/99-own-defaults
uci -q set irqbalance.enabled='1'
uci commit irqbalance

uci -q set globals.packet_steering ='1'
uci commit globals
